# deploy/systemd/bizscout-run.timer
# =============================================================
#  Bizreach Auto Scout — 初回送信タイマー（毎日）
# -------------------------------------------------------------
#  bizscout-run.service を毎日決まった時刻に起動する。
# =============================================================
[Unit]
Description=Bizreach Auto Scout - 初回送信を毎日実行するタイマー

[Timer]
# 毎日 09:00 に実行（OnCalendar=曜日 年-月-日 時:分:秒 の書式。曜日省略で毎日）。
OnCalendar=*-*-* 09:00:00
# マシン停止等で実行を逃した場合、起動後に取り戻して実行する。
Persistent=true
# 同時刻集中を避けるための最大ランダム遅延（任意）。
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
