# deploy/systemd/bizscout-resends.timer
# =============================================================
#  Bizreach Auto Scout — 再送タイマー（毎日）
# -------------------------------------------------------------
#  bizscout-resends.service を毎日決まった時刻に起動する。
#  初回送信ジョブ(bizscout-run)と時刻をずらして運用するとよい。
# =============================================================
[Unit]
Description=Bizreach Auto Scout - 再送を毎日実行するタイマー

[Timer]
# 毎日 10:00 に実行（初回送信の 09:00 とずらす）。
OnCalendar=*-*-* 10:00:00
# マシン停止等で実行を逃した場合、起動後に取り戻して実行する。
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
